# build/supervisor.conf
[program:app-django]
command = python3 /home/ubuntu/src/manage.py runserver 0.0.0.0:8000